# tests/material/supervisord/supervisord.conf
[supervisord]

logfile=/dev/null
logfile_maxbytes=0

[program:linear_algebra]

command=/bin/bash -c "python3 -m superduperdb.servers.linear_algebra $(cat config.json | jq .linear_algebra.port)"
process_name=%(program_name)s_%(process_num)s
numprocs=1
stdout_logfile=logs/linear_algebra.out
stderr_logfile=logs/linear_algebra.out
autorestart=false
startretries=1

[program:model_server]

command=/bin/bash -c "python3 -m superduperdb.servers.model_server $(cat config.json | jq .model_server.port)"
process_name=%(program_name)s_%(process_num)s
numprocs=1
stdout_logfile=logs/model_server.out
stderr_logfile=logs/model_server.out
autorestart=false
startretries=1

[program:worker]

command=/bin/bash -c "rq worker -v --url redis://:@localhost:$(cat config.json | jq .redis.port)"
process_name=%(program_name)s_%(process_num)s
numprocs=2
stdout_logfile=logs/worker.out
stderr_logfile=logs/worker.out
autorestart=false
startretries=1

[program:redis]

command=/bin/bash -c "redis-server --port $(cat config.json | jq .redis.port)"
process_name=%(program_name)s_%(process_num)s
numprocs=1
stdout_logfile=logs/redis.out
stderr_logfile=logs/redis.out
autorestart=false
startretries=1

[program:image_server]

command=/bin/bash -c "python3 -m tests.material.supervisord.image_server"
process_name=%(program_name)s_%(process_num)s
numprocs=1
stdout_logfile=logs/image_server.out
stderr_logfile=logs/image_server.out
autorestart=false
startretries=1
