# tests/material/supervisord/supervisord_rq.conf
[supervisord]

logfile=logs/supervisord.log

[program:vector_search]

command=/bin/bash -c "python3 -m superduperdb.serving.vector_search $(cat config.json | jq .linear_algebra.port)"
process_name=%(program_name)s_%(process_num)s
numprocs=1
stdout_logfile=logs/vector_search.out
stderr_logfile=logs/vector_search.out
autorestart=false
startretries=1

[program:model_server]

command=/bin/bash -c "python3 -m superduperdb.serving.model_server $(cat config.json | jq .model_server.port)"
process_name=%(program_name)s_%(process_num)s
numprocs=1
stdout_logfile=logs/model_server.out
stderr_logfile=logs/model_server.out
autorestart=false
startretries=1

[program:worker]

command=/bin/bash -c "rq worker -w rq_worker.py -v --url redis://:@localhost:$(cat config.json | jq .redis.port)"
process_name=%(program_name)s_%(process_num)s
numprocs=1
stdout_logfile=logs/worker_%(process_num)s.out
stderr_logfile=logs/worker_%(process_num)s.out
autorestart=false
startretries=1

[program:redis]

command=/bin/bash -c "redis-server --port $(cat config.json | jq .redis.port)"
process_name=%(program_name)s_%(process_num)s
numprocs=1
stdout_logfile=logs/redis.out
stderr_logfile=logs/redis.out
autorestart=false
startretries=1

