# tests/material/supervisord/supervisord_dask.conf
[supervisord]

logfile=logs/supervisord.log

[program:vector_search]

command=/bin/bash -c "python3 -m superduperdb.cluster.vector_search $(cat config.json | jq .linear_algebra.port)"
process_name=%(program_name)s_%(process_num)s
numprocs=1
stdout_logfile=logs/vector_search.out
stderr_logfile=logs/vector_search.out
autorestart=false
startretries=1

[program:model_server]

command=/bin/bash -c "python3 -m superduperdb.cluster.model_server $(cat config.json | jq .model_server.port)"
process_name=%(program_name)s_%(process_num)s
numprocs=1
stdout_logfile=logs/model_server.out
stderr_logfile=logs/model_server.out
autorestart=false
startretries=1

[program:dask_worker]

command=/bin/bash -c "dask-worker tcp://127.0.0.1:8786"
process_name=%(program_name)s_%(process_num)s
numprocs=1
stdout_logfile=logs/dask_worker_%(process_num)s.out
stderr_logfile=logs/dask_worker_%(process_num)s.out
autorestart=false
startretries=1

[program:dask_scheduler]

command=/bin/bash -c "dask-scheduler --host 127.0.0.1"
process_name=%(program_name)s_%(process_num)s
numprocs=1
stdout_logfile=logs/dask_scheduler.out
stderr_logfile=logs/dask_scheduler.out
autorestart=false
startretries=1

